# deploy/systemd/treesql-console.service
[Unit]
Description=TreeSQL Console

[Service]
ExecStart=/home/ubuntu/ran -r /home/ubuntu/webui -p 80

[Install]
WantedBy=multi-user.target
