# deploy/systemd/treesql-server.service
[Unit]
Description=TreeSQL server

[Service]
ExecStart=/home/ubuntu/treesql-server

[Install]
WantedBy=multi-user.target
