# deploy/systemd/treesql-console.service
[Unit]
Description=TreeSQL Console

[Service]
ExecStart=/usr/bin/ran -r /usr/bin/treesql-webui -p 80

[Install]
WantedBy=multi-user.target
